# monad-solonet/solonet/services/monad-execution.ini
[program:monad-execution]
autostart=false
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile=/var/log/monad-execution.log
stderr_logfile=/var/log/monad-execution.log
command=run-monad-execution.sh
