# monad-solonet/solonet/services/monad-bft.ini
[program:monad-bft]
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile=/var/log/monad-bft.log
stderr_logfile=/var/log/monad-bft.log
environment=
  RUST_LOG="debug,h2=warn,tower=warn,opentelemetry_sdk=warn,opentelemetry-otlp=warn",
  REMOTE_FORKPOINT_URL=http://localhost:8082/shared/forkpoint.toml,
  REMOTE_VALIDATORS_URL=http://localhost:8082/shared/validators.toml
command=cpulimit --foreground -l 50 -- monad-node
  --triedb-path=/dev/triedb
  --secp-identity=/home/monad/monad-bft/config/id-secp
  --bls-identity=/home/monad/monad-bft/config/id-bls
  --node-config=/home/monad/monad-bft/config/node.toml
  --forkpoint-config=/home/monad/monad-bft/config/forkpoint/forkpoint.toml
  --wal-path=/home/monad/monad-bft/wal
  --mempool-ipc-path=/home/monad/monad-bft/mempool.sock
  --persisted-peers-path=/home/monad/monad-bft/config/peers.toml
  --control-panel-ipc-path=/home/monad/monad-bft/controlpanel.sock
  --statesync-ipc-path=/home/monad/monad-bft/statesync.sock
  --ledger-path=/home/monad/monad-bft/ledger
  --otel-endpoint="http://0.0.0.0:4317"
  --record-metrics-interval-seconds=1
  --validators-path=/home/monad/monad-bft/config/validators/validators.toml
  --keystore-password=password
