# monad-solonet/solonet/services/monad-ledger-tail.ini
[program:monad-ledger-tail]
autostart=false
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile=/var/log/monad-ledger-tail.log
stderr_logfile=/var/log/monad-ledger-tail.log
command=run-monad-ledger-tail.sh
