# monad-solonet/solonet/services/monad-txgen.ini
[program:monad-txgen]
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile=/var/log/monad-txgen.log
stderr_logfile=/var/log/monad-txgen.log
environment=
  RUST_LOG="debug,h2=warn,tower=warn,opentelemetry_sdk=warn,opentelemetry-otlp=warn"
command=monad-txgen --config-file /solonet/config/txgen.toml
