# monad-solonet/solonet/services/monad-txgen.ini
[program:monad-txgen]
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile=/var/log/monad-txgen.log
stderr_logfile=/var/log/monad-txgen.log
command=run-monad-txgen.sh
