# monad-solonet/solonet/services/monad-execution.ini
[program:monad-execution]
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile=/var/log/monad-execution.log
stderr_logfile=/var/log/monad-execution.log
command=cpulimit --foreground -l 50 -- monad
  --chain=monad_devnet
  --db=/dev/triedb
  --block_db=/home/monad/monad-bft/ledger
  --statesync=/home/monad/monad-bft/statesync.sock
  --log_level=ERROR
