# monad-solonet/solonet/services/monad-ledger-tail.ini
[program:monad-ledger-tail]
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile=/var/log/monad-ledger-tail.log
stderr_logfile=/var/log/monad-ledger-tail.log
command=monad-ledger-tail
  --ledger-path=/home/monad/monad-bft/ledger
  --forkpoint-path=/home/monad/monad-bft/config/forkpoint/forkpoint.toml
  --node-config-path=/home/monad/monad-bft/config/node.toml
  --validators-path=/home/monad/monad-bft/config/validators/validators.toml
