# monad-solonet/solonet/services/monad-bft.ini
[program:monad-bft]
autostart=false
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile=/var/log/monad-bft.log
stderr_logfile=/var/log/monad-bft.log
command=/solonet/services/start-monad-bft.sh
