# monad-solonet/solonet/services/nginx.ini
[program:nginx]
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.log
command=/usr/sbin/nginx -g "daemon off;"
