# monad-solonet/solonet/services/sync-forkpoint-files.ini
[program:sync-forkpoint-files]
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile_maxbytes = 10MB
stderr_logfile_maxbytes = 10MB
stdout_logfile=/var/log/sync-forkpoint-files.log
stderr_logfile=/var/log/sync-forkpoint-files.log
command=sync-forkpoint-files.sh
